# nflabs-status-agent/crontab
# run every 5 minutes

*/5 * * * * /path/to/nflabs_status_agent.sh
